# image_file/rc.local
#!/bin/sh
echo "common operaions..."

echo "sync time..."
ntpdate time.google.com

echo "mounting gcs..."
gcsfuse serval-data /app/data/gcloud